# conf/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear -a al %I $TERM
#ExecStart=-/sbin/agetty --autologin al --noclear %I $TERM
