# conf/etc/systemd/system/mnt-remote.automount
[Unit]
Description=Automount NFS share after VPN is up

[Automount]
Where=/mnt/remote

[Install]
WantedBy=multi-user.target
