# conf/etc/systemd/system/wg-quick@ai.service.d/override.conf
[Unit]
Requires=wstunnel-client.service
