# conf/home/config/autostart/gnome-terminal.desktop
[Desktop Entry]
Type=Application
Name=gnome-terminal
Exec=/usr/bin/gnome-terminal --full-screen
