# conf/etc/systemd/system/wg-quick@z.service.d/override.conf
[Service]
ExecStartPre=bash -c "/usr/bin/sed -i /Endpoint/s:..$:$(date +%%d): /etc/wireguard/%i.conf"
